# files/default/init.d/serf
#!/bin/sh
#
# serf Starts/stop the "at" daemon
#
# chkconfig:   345 95 5
#

# Source function library.
. /etc/rc.d/init.d/functions

exec="/opt/serf/sbin/serf"
logfile="/opt/serf/var/log/serf.log"
prog="serf"
config="/opt/serf/etc/serf.json"
options="agent -config-file=$config"

start() {
    [ -x $exec ] || exit 5
    [ -f $config ] || exit 6
    echo $"Starting ${prog}... "
    $exec $options >> $logfile &
    sleep 1
    echo "Done."
}

stop() {
    echo $"Stopping ${prog}... "
    $exec leave
    echo "Done."
}

restart() {
    stop
    start
}

case "$1" in
    start)
        $1
        ;;
    stop)
        $1
        ;;
    restart)
        $1
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart}"
        exit 2
esac
exit $?
